# rpm/SOURCES/shadowsocks-server@.service
[Unit]
Description=Shadowsocks Server Service
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/ss-server -c /etc/shadowsocks/%i.json

[Install]
WantedBy=multi-user.target
