# openwrt/files/shadowsocks.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=95
STOP=95

start() {
    echo "starting shadowsocks..."
    /usr/bin/ss-local -c /etc/config/shadowsocks.json -f /var/run/shadowsocks.pid
}

stop() {
    echo "shadowsocks stopped."
    kill -9 `cat /var/run/shadowsocks.pid`
}

